# deploy/moltbot-gateway.service
[Unit]
Description=Moltbot Gateway - Personal AI Assistant
Documentation=https://docs.molt.bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=moltbot
Group=moltbot
WorkingDirectory=/home/moltbot
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=1536
Environment=PATH=/home/moltbot/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/moltbot
EnvironmentFile=-/home/moltbot/.config/moltbot/.env
ExecStartPre=+/bin/sh -c 'mkdir -p /home/moltbot/.clawdbot /home/moltbot/clawd/memory && chown -R moltbot:moltbot /home/moltbot/.clawdbot /home/moltbot/clawd && chmod 700 /home/moltbot/.clawdbot /home/moltbot/clawd'
ExecStartPre=/bin/sh -c 'echo "moltbot-gateway: pre-start checks..." && test -x /home/moltbot/.npm-global/bin/moltbot || { echo "FATAL: /home/moltbot/.npm-global/bin/moltbot not found or not executable"; exit 1; } && test -f /home/moltbot/.config/moltbot/.env || echo "WARN: /home/moltbot/.config/moltbot/.env not found, running without env file"'
ExecStart=/home/moltbot/.npm-global/bin/moltbot gateway --port 18789
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moltbot

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/moltbot/.config/moltbot /home/moltbot/.local/share/moltbot /home/moltbot/.npm-global /home/moltbot/.clawdbot /home/moltbot/clawd
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# Resource limits
# NOTE: install.sh and update.sh auto-tune MemoryMax and --max-old-space-size
# based on detected system memory (see compute_memory_limits in lib.sh).
# Formula: heap = 65% RAM (floor 256M, cap 1536M)
#          MemoryMax = heap + 512M overhead (floor 512M, cap 2G, ≤ 90% RAM)
# The values below are static defaults for 4 GB+ systems.
LimitNOFILE=65535
MemoryMax=2G

[Install]
WantedBy=multi-user.target
