# deploy/moltbot-gateway.service
[Unit]
Description=Moltbot Gateway - Personal AI Assistant
Documentation=https://docs.molt.bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=moltbot
Group=moltbot
WorkingDirectory=/home/moltbot
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=1536
Environment=PATH=/home/moltbot/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/moltbot
EnvironmentFile=-/home/moltbot/.config/moltbot/.env
ExecStartPre=/bin/sh -c 'echo "moltbot-gateway: pre-start checks..." && test -x /home/moltbot/.npm-global/bin/moltbot || { echo "FATAL: /home/moltbot/.npm-global/bin/moltbot not found or not executable"; exit 1; } && test -f /home/moltbot/.config/moltbot/.env || echo "WARN: /home/moltbot/.config/moltbot/.env not found, running without env file"'
ExecStart=/home/moltbot/.npm-global/bin/moltbot gateway --port 18789
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moltbot

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/moltbot/.config/moltbot /home/moltbot/.local/share/moltbot /home/moltbot/.npm-global
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# Resource limits
# NOTE: install.sh auto-tunes MemoryMax (75% of RAM, max 2G) and
# NODE_OPTIONS --max-old-space-size (50% of RAM, max 1536M) based on
# detected system memory. The values below are defaults for 4GB+ systems.
LimitNOFILE=65535
MemoryMax=2G

[Install]
WantedBy=multi-user.target
